# root/etc/services.d/booksonic/run
#!/usr/bin/with-contenv bash

cd /app/booksonic
exec \
	s6-setuidgid abc java -jar /usr/share/java/jetty-runner.jar \
	/app/booksonic/booksonic.war
